# lcd_display.path
[Unit]
Description=Watch for network changes and update LCD

[Path]
PathChanged=/run/NetworkManager/
PathChanged=/sys/class/net/
Unit=lcd_display.service

[Install]
WantedBy=multi-user.target
